# deploy/systemd/arv.service
[Unit]
Description=Vertex AR FastAPI Application
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=exec
User=arv
Group=arv
WorkingDirectory=/opt/arv/app
EnvironmentFile=/opt/arv/app/.env

ExecStart=/opt/arv/venv/bin/gunicorn app.main:app \
    --bind 127.0.0.1:8000 \
    --workers 2 \
    --worker-class uvicorn.workers.UvicornWorker \
    --timeout 120 \
    --graceful-timeout 30 \
    --keep-alive 5 \
    --max-requests 1000 \
    --max-requests-jitter 50 \
    --access-logfile - \
    --log-level info

Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=arv

[Install]
WantedBy=multi-user.target
